# doc/config/nginx_supervisor.conf
[program:nginx]
directory=/usr/local/nginx ;#项目目录
user=root
command=/usr/local/nginx/sbin/nginx ;#启动命令
autostart=true
autorestart=true
stdout_logfile=/tmp/nginx.log ;#日志
redirect_stderr=true